# systemd/media-stack.service
[Unit]
Description=Media Stack Docker Compose Service (Audiobookshelf & Navidrome)
Documentation=https://docs.docker.com/compose/
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/homelab
ExecStart=/usr/bin/docker compose -f docker-compose.media.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.media.yml down
TimeoutStartSec=0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
