# systemd/storagebox-mount.service
[Unit]
Description=Mount Hetzner Storage Box via WebDAV
Documentation=https://docs.hetzner.com/storage/storage-box/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/mount-storagebox.sh
ExecStop=/bin/umount /mnt/storagebox
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

